# scripts/mediamtx/com.mediamtx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mediamtx</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/mediamtx</string>
        <string>/Users/alpaca/mediamtx/mediamtx.yml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/alpaca/mediamtx</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/alpaca/mediamtx/mediamtx.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alpaca/mediamtx/mediamtx.log</string>
</dict>
</plist>
