# bug-fixer/bug-fixer.service
[Unit]
Description=GenAlpaca Bug Scout
After=network.target

[Service]
Type=simple
User=bugfixer
Environment=HOME=/home/bugfixer
WorkingDirectory=/opt/bug-fixer
EnvironmentFile=/opt/bug-fixer/.env
ExecStart=/usr/bin/node /opt/bug-fixer/bug_scout.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bug-fixer

[Install]
WantedBy=multi-user.target
