# tesla-poller/tesla-poller.service
[Unit]
Description=GenAlpaca Tesla Vehicle Data Poller
After=network.target

[Service]
Type=simple
User=bugfixer
Environment=HOME=/home/bugfixer
WorkingDirectory=/opt/tesla-poller
EnvironmentFile=/opt/tesla-poller/.env
ExecStart=/usr/bin/node /opt/tesla-poller/worker.js
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tesla-poller

[Install]
WantedBy=multi-user.target
